# scripts/food-ordering.service
# ================================================================
# فایل سرویس سیستم برای سیستم سفارش غذا
# این فایل برای راه‌اندازی خودکار سیستم به عنوان سرویس سیستم‌عامل استفاده می‌شود
# نویسنده: تیم توسعه
# تاریخ آخرین ویرایش: تیر ۱۴۰۴
# نسخه: ۲.۰ - سرویس پیشرفته با مدیریت کامل
# ================================================================

[Unit]
# ================================================================
# بخش تنظیمات واحد سرویس
# ================================================================
Description=Food Ordering System Service                    # توضیحات سرویس
Documentation=https://github.com/food-ordering-system      # مستندات پروژه
After=network.target                                       # اجرا پس از راه‌اندازی شبکه
Wants=network.target                                       # نیاز به شبکه
Requires=network.target                                    # اجباری بودن شبکه

# ================================================================
# بخش تنظیمات سرویس
# ================================================================
[Service]
# نوع سرویس: ساده (simple) برای برنامه‌های ساده
Type=simple

# کاربر و گروه اجراکننده سرویس
User=food-ordering                                         # نام کاربر
Group=food-ordering                                        # نام گروه

# مسیر کاری سرویس
WorkingDirectory=/opt/food-ordering-system                 # مسیر کاری

# دستور اجرای سرویس
ExecStart=/usr/bin/java -jar /opt/food-ordering-system/backend/target/food-ordering-backend.jar  # دستور شروع

# دستور توقف سرویس
ExecStop=/bin/kill -TERM $MAINPID                         # دستور توقف

# دستور بارگذاری مجدد سرویس
ExecReload=/bin/kill -HUP $MAINPID                        # دستور بارگذاری مجدد

# تنظیمات محیطی
Environment=JAVA_HOME=/usr/lib/jvm/java-17-openjdk        # مسیر Java
Environment=SPRING_PROFILES_ACTIVE=production             # پروفایل فعال Spring
Environment=DB_PATH=/opt/food-ordering-system/backend/food_ordering.db  # مسیر پایگاه داده

# تنظیمات امنیتی
NoNewPrivileges=true                                       # عدم دسترسی به امتیازات جدید
PrivateTmp=true                                            # استفاده از /tmp خصوصی
ProtectSystem=strict                                        # محافظت شدید از سیستم
ProtectHome=true                                           # محافظت از پوشه home
ReadWritePaths=/opt/food-ordering-system                   # مسیرهای قابل نوشتن

# تنظیمات منابع
LimitNOFILE=65536                                          # حداکثر تعداد فایل‌های باز
LimitNPROC=4096                                            # حداکثر تعداد پردازش‌ها
MemoryMax=2G                                               # حداکثر حافظه استفاده شده

# تنظیمات راه‌اندازی مجدد
Restart=always                                             # همیشه راه‌اندازی مجدد
RestartSec=10                                              # فاصله زمانی راه‌اندازی مجدد (ثانیه)
StartLimitInterval=60                                      # فاصله زمانی محدودیت شروع (ثانیه)
StartLimitBurst=3                                          # حداکثر تعداد شروع در فاصله زمانی

# تنظیمات تایم‌اوت
TimeoutStartSec=300                                        # تایم‌اوت شروع (ثانیه)
TimeoutStopSec=60                                          # تایم‌اوت توقف (ثانیه)

# تنظیمات لاگینگ
StandardOutput=journal                                     # خروجی استاندارد به journal
StandardError=journal                                      # خروجی خطا به journal
SyslogIdentifier=food-ordering                             # شناسه syslog

# ================================================================
# بخش تنظیمات نصب
# ================================================================
[Install]
# سرویس‌های مورد نیاز برای نصب
WantedBy=multi-user.target                                 # نصب در حالت چندکاربره

# ================================================================
# بخش تنظیمات اضافی (اختیاری)
# ================================================================
# تنظیمات اضافی برای بهینه‌سازی عملکرد
# KillMode=mixed                                           # حالت توقف پردازش‌ها
# KillSignal=SIGTERM                                       # سیگنال توقف
# SendSIGKILL=yes                                          # ارسال SIGKILL در صورت عدم توقف
# SendSIGHUP=no                                            # عدم ارسال SIGHUP 